# aur-repo/csnet/csnet.service
[Unit]
Description=csnet service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=csnetuser
Group=csnetuser
WorkingDirectory=/var/lib/csnet
ExecStart=/usr/bin/csnet -c /etc/csnet/client_produce.json
PrivateTmp=true
Restart=always
RestartSec=30s
StandardOutput=append:/var/log/csnet/csnet.log
StandardError=append:/var/log/csnet/csnet.log

[Install]
WantedBy=multi-user.target
