# aur-repo/openlist/openlist.service
[Unit]
Description=openlist
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=openlist
WorkingDirectory=/var/lib/openlist
ExecStart=/usr/bin/openlist server --data /var/lib/openlist
Restart=on-failure

[Install]
WantedBy=multi-user.target
