# aur-repo/iyuuplus/iyuuplus.service
[Unit]
Description = IYUU Auto Reseed Plus Daemon.
Documentation = https://doc.iyuu.cn
Wants = network-online.target mariadb.service
After = network-online.target mariadb.service

[Service]
User = iyuuplus
Group = iyuuplus
Restart = on-failure
ExecStart = /usr/bin/iyuuplus start

[Install]
WantedBy = multi-user.target
