# aur-repo/ddnsto-bin/ddnsto.service
[Unit]
Description=DDNSTO Service
Documentation=https://doc.linkease.com/zh/guide/ddnsto
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
Type=simple
User=ddnsto
Group=ddnsto

PrivateTmp=yes
ProtectSystem=no
ProtectHome=no

ExecStart=/usr/bin/ddnsto --config /var/lib/ddnsto/config.yaml
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=ddnsto

[Install]
WantedBy=multi-user.target
