# aur-repo/ddnsto-bin/ddnsto-bin.install
_help() {
    cat << 'EOF'
================================
DDNSTO 快速使用指南 / Quick Start Guide
================================

1. 配置并启动服务 / Configure and Start Service
sudo vim /etc/ddnsto/config.yaml
# 必须修改: userToken routerName 和 webdav.password
# Must change: userToken routerName and webdav.password
sudo systemctl enable --now ddnsto

2. 检查服务状态 / Check Service Status
sudo systemctl status ddnsto

3. 查看实时日志 / View Real-time Logs
sudo journalctl -u ddnsto -f

4. WebDAV 访问 / WebDAV Access
账号/Username: admin
密码/Password: config.yaml 中设置的密码
密码/Password: Set in config.yaml
路径/Path: /srv/ddnsto/
在浏览器打开 / Open in browser:
xdg-open http://localhost:18088

================================
EOF
}

post_install() {
   _help
}

post_upgrade() {
  _help
}