# aur-repo/ddnsto-bin/ddnsto-bin.install
_help() {
    cat << 'EOF'
================================
DDNSTO 快速使用
================================

1. 配置
sudo vim /etc/ddnsto/config.yaml
# 修改: userToken routerName 和 webdav.password

2. 启动
sudo systemctl enable --now ddnsto

3. 检查状态
sudo systemctl status ddnsto

4. 查看日志
sudo journalctl -u ddnsto -f

5. 访问 WebDAV
浏览器打开: http://localhost:18088
账号: admin
密码: config.yaml 中设置
路径: /srv/ddnsto/
================================
EOF
}

post_install() {
   _help
}

post_upgrade() {
  _help
}