# aur-repo/sunloginclient/runsunloginclient.service
[Unit]
Description=Sunlogin (AweSun) remote control daemon
Requires=network-online.target
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/awesun/bin/awesun_daemon -m server -name awesun
KillMode=control-group
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
