# aur-repo/wyc/wyc.service
[Unit]
Description=WangYunChuan
After=network.target

[Service]
Type=simple
EnvironmentFile=/etc/wyc.conf
ExecStart=/usr/bin/wyc -token=${TOKEN}
PrivateTmp=true
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=multi-user.target
